# antbot.service
[Unit]
Description=AntBot

[Service]
ExecStart=home/AntBot/startup.sh

[Install]
WantedBy=multi-user.target